# contrib/linux-packaging/archlinux/neucoind.service
[Unit]
Description=Bfxcoin daemon
After=network.target

[Service]
ExecStart=bfxcoind -logtimestamps

[Install]
WantedBy=default.target
